# supervisord.conf

[supervisord]
nodaemon=true


[program:influx]
command=/bin/bash -c "exec /usr/bin/influxd -config /etc/influxdb/influxdb.conf"

[program:database]
command=/bin/bash -c "exec ./run.sh"
